# .config/autostart/01-qv2ray.desktop
[Desktop Entry]
Name=Qv2ray Delayed
Exec=/bin/sh -c "sleep 1s; cgnoproxy qv2ray"
Terminal=false
Type=Application
StartupNotify=false
