# .config/autostart/01-qv2ray.desktop
[Desktop Entry]
Name=Qv2ray (Delayed 1s)
Exec=bash -c "sleep 1s; qv2ray"
Terminal=false
Type=Application
StartupNotify=false
