# .config/autostart/02-flameshot.desktop
[Desktop Entry]
Name=Flameshot (Delayed 2s)
Exec=bash -c "sleep 2s; flameshot"
Terminal=false
Type=Application
StartupNotify=false
