# .config/autostart/03-crow-translate.desktop
[Desktop Entry]
Name=Crow Translate (Delayed 3s)
Exec=bash -c "sleep 3s; crow"
Terminal=false
Type=Application
StartupNotify=false
